# DEBIAN/postinst
#!/bin/bash
set -e

systemctl enable radxa-e25-led.service || true
